# deploy/systemd/mimo-code-webui.service
[Unit]
Description=MiMo Code WebUI
Documentation=https://github.com/yaozy2020/mimo-code-webui
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mimo-webui
Group=mimo-webui
WorkingDirectory=/opt/mimo-code-webui/current
EnvironmentFile=/etc/mimo-code-webui/webui.env
ExecStart=/opt/mimo-code-webui/current/scripts/start.sh
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStartSec=60
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
