# deploy/systemd/mimo-code-webui-backup.timer
[Unit]
Description=Daily MiMo Code WebUI backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
