# deploy/systemd/mimo-code-webui-backup.service
[Unit]
Description=Back up MiMo Code WebUI and MiMo state
Documentation=https://github.com/yaozy2020/mimo-code-webui
ConditionPathExists=/opt/mimo-code-webui/current/scripts/backup-state.mjs
OnFailure=mimo-code-webui-alert@%n.service

[Service]
Type=oneshot
EnvironmentFile=/etc/mimo-code-webui/webui.env
Environment=MIMO_BACKUP_OFFLINE=true
ExecStart=/usr/bin/env MIMO_BACKUP_OPERATION_ID=${INVOCATION_ID} /opt/mimo-code-webui/current/scripts/backup-service-state.sh run
ExecStopPost=/usr/bin/env MIMO_BACKUP_OPERATION_ID=${INVOCATION_ID} /opt/mimo-code-webui/current/scripts/backup-service-state.sh restore
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
